# cgiproxy.cron
13 2 * * * root /opt/cgiproxy/bin/nph-proxy.cgi purge-db
