# cgiproxy.cron
13 2 * * * /opt/cgiproxy/bin/nph-proxy.cgi purge-db
