# services/cgiproxy/run
#!/bin/bash
exec setuser proxy /cgiproxy/bin/nph-proxy.cgi start-fcgi
